# mavlink-timesync.service
[Unit]
Description=MAVLink System Time Synchronization
After=network.target

[Service]
ExecStart=/usr/local/bin/mavlinktimesync.py
User=root
Restart=always

[Install]
WantedBy=multi-user.target
